# supervisord.conf
[program:init]
autorestart=true
command=/usr/bin/env python3 /tmp/init.py

[program:rundeck]
autorestart=false
startsecs=0
command=service rundeckd start
